# samba/post-install/etc/service/avahi/run
#!/bin/bash -x

exec 2>&1

set -eu -o pipefail

sv check /etc/service/dbus

exec /usr/sbin/avahi-daemon --debug
