# samba/post-install/etc/service/dbus/run
#!/bin/bash

exec 2>&1

[[ "$DEBUG" -ge 1 ]] && set -x

set -eu -o pipefail

mkdir -p /run/dbus
rm -vf /run/dbus.pid

dbus-uuidgen --ensure

exec dbus-daemon --system --nosyslog --nofork
