# samba/post-install/etc/service/samba/run
#!/bin/bash

exec 2>&1

[[ "$DEBUG" -ge 1 ]] && set -x

set -eu -o pipefail

sv check /etc/service/avahi

testparm -s

inspect-xattr

exec /usr/sbin/smbd --log-stdout --foreground
